# crontab
*/5 * * * * $HOME/temps/collect-temp-cronjob.sh
